# deploy/rag.service
[Unit]
Description=Notes RAG — hybrid search + file watcher
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rag
ExecStart=/opt/rag/venv/bin/python /opt/rag/main.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
# Uncomment and set if using Vault:
# Environment=VAULT_ADDR=http://192.168.88.65:8200
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
